# 40-rapoo-keyboard-driver.rules
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="1c4f", ATTR{idProduct}=="0059", RUN+="/usr/local/bin/install-rapoo-keyboard-driver.sh"
